# systemd/cpj.service
[Unit]
Description=CPJ service
Documentation= man:systemd.service

[Service]
Type=oneshot
ExecStart=/opt/cpj/bin/cpj.sh
