# systemd/cpj.service
[Unit]
Description=Container PKI Janitor service
Documentation=man:systemd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cpj scan
