# systemd/cpj.timer
[Unit]
Description=Container PKI Janitor timer
Documentation=man:systemd.timer

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=multi-user.target
