# systemd/cpj.path
[Unit]
Description=Container PKI Janitor path monitor
Documentation=man:systemd.path

[Path]
PathModified=/var/lib/docker/containers

[Install]
WantedBy=multi-user.target
